# supervisord.conf
[supervisord]
nodaemon=true

[program:puppetmaster]
command=/bin/bash -c "puppet master --verbose --no-daemonize"
